# scripts/agent-hunter.timer
[Unit]
Description=agent-hunter 每日运行定时器

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
